# roles/base/templates/etc/profile.d/server-env.sh.j2
#!/bin/sh
#
# Server Environment Configuration
#

export APPLICATION_ENV={{server_env}}
